# templates/lcd.service
[Unit]
Description=LCD Service

[Service]
ExecStart=/usr/bin/lcd
Restart=always
User=root

[Install]
WantedBy=multi-user.target